# goodies/StarinetPythonLogger.conf
[program:StarinetPythonLogger]
directory=/opt/StarinetPythonLogger
command=python main.py
autostart=true
autorestart=true
