# goodies/StarinetPythonLogger.conf
[program:StarinetPythonLogger]
directory=/opt/StarinetPythonLogger
command=python StarinetPython3Logger.py
autostart=true
autorestart=true
